# config/launchagents/com.aos.qareen-tunnel.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aos.qareen-tunnel</string>
    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/.aos/services/qareen/.venv/bin/python</string>
        <string>-c</string>
        <string>import os,subprocess; subprocess.run(["__CLOUDFLARED__","tunnel","run","--token",os.environ["CF_TUNNEL_TOKEN"]])</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>5</integer>
    <key>StandardOutPath</key>
    <string>__HOME__/.aos/logs/qareen-tunnel.out.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/.aos/logs/qareen-tunnel.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:__HOME__/.bun/bin:__HOME__/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>CF_TUNNEL_TOKEN</key>
        <string>__CF_TUNNEL_TOKEN__</string>
    </dict>
</dict>
</plist>
